# etc/systemd/system/rotate-vpn.timer
[Unit]
Description=Rotates the autoconnect VPN every midnight

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target